# local/dashboard/dashboard.service
[Unit]
Description=TODO-HOSTNAME Dashboard (Bun)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sysadmin-agent/local/dashboard
EnvironmentFile=/home/ubuntu/sysadmin-agent/local/.env
Environment=PATH=/home/ubuntu/.bun/bin:/snap/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/.bun/bin/bun run server.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
