# templates/local/dashboard/dashboard.service
# Sysadmin Dashboard — systemd unit template
#
# Placeholders (substituted by setup.sh or manually):
#   __USER__      — Linux user running the agent
#   __REPO__      — Absolute path to the sysadmin-agent repo
#   __HOME__      — Home directory of __USER__
#   __HOSTNAME__  — Machine hostname (cosmetic)
#
# Manual install:
#   sed -e 's|__USER__|ubuntu|g' \
#       -e 's|__REPO__|/home/ubuntu/sysadmin-agent|g' \
#       -e 's|__HOME__|/home/ubuntu|g' \
#       -e 's|__HOSTNAME__|myhost|g' \
#       templates/local/dashboard/dashboard.service \
#       | sudo tee /etc/systemd/system/myhost-dashboard.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now myhost-dashboard

[Unit]
Description=__HOSTNAME__ Dashboard (Bun)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO__/local/dashboard
EnvironmentFile=__REPO__/local/.env
Environment=PATH=__HOME__/.bun/bin:/snap/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__HOME__/.bun/bin/bun run server.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
