# scripts/telegram-bot/sysadmin-agent-telegram.service
# DEPRECATED — Use interactive Claude session with MCP plugin instead.
# Kept as reference only. If needed, install via setup.sh or manually substitute
# %REPO_PATH% and %USER% before copying to /etc/systemd/system/.

[Unit]
Description=Sysadmin Agent Telegram Dispatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%REPO_PATH%
ExecStart=%REPO_PATH%/scripts/telegram-bot/telegram-dispatch.sh
Restart=on-failure
RestartSec=30
EnvironmentFile=%REPO_PATH%/local/.env

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
